# dev-db/couchbase-server-bin/files/couchbase-server
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/ntp/files/ntpd.rc,v 1.24 2010/07/18 21:53:24 vapier Exp $

BASEDIR=/opt/couchbase/var/lib/couchbase
DAEMON=/opt/couchbase/bin/couchbase-server
PIDFILE=/opt/couchbase/var/lib/couchbase/couchbase-server.pid
NODEFILE=/opt/couchbase/var/lib/couchbase/couchbase-server.node
COOKIEFILE=/opt/couchbase/var/lib/couchbase/couchbase-server.cookie

depend() {
	use net
}

checkconfig() {
	ulimit -n 10240
	ulimit -c unlimited
	ulimit -l unlimited

	touch ${PIDFILE} ${NODEFILE} ${COOKIEFILE} && chown couchbase:couchbase ${PIDFILE} ${NODEFILE} ${COOKIEFILE}
	return $?
}

start() {
	checkconfig || return $?

	ebegin "Starting ${SVCNAME}"
	start-stop-daemon --start --chdir ${BASEDIR} --user couchbase --pidfile ${PIDFILE} \
		--exec /bin/sh -- -c "${DAEMON} -- -noinput -detached > /opt/couchbase/var/lib/couchbase/logs/start.log 2>&1"

	einfon "Waiting for startup ."
	nc=0
	while true; do
		[ -f ${PIDFILE} ] && ps --pid $(cat ${PIDFILE}) &>/dev/null && echo && break
		nc=$(expr $nc + 1)
		if [ $nc -gt 60 ]; then
			echo
			eend 1 "Timeout waiting for ${SVCNAME} startup !"
			exit 1
		fi
		sleep 1
		echo -n "."
	done

	eend $? "Failed to start ${SVCNAME}"
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	${DAEMON} -k
	eend $? "Failed to stop ${SVCNAME}"
}
